# rootfs/etc/s6/pdns_server/run
#!/bin/sh

exec /usr/sbin/pdns_server --launch=bind,lua2 --lua2-filename=/etc/pdns/powerdns-luabackend.lua --daemon=no --guardian=no --disable-syslog
