# rootfs/etc/s6/pdns_recursor/run
#!/bin/sh

mkdir -p /var/run/pdns-recursor
exec /usr/sbin/pdns_recursor --daemon=no --disable-syslog
